# client/rc.local
#!/bin/bash

rm /usr/local/share/ca-certificates/fakeroot.crt
while [ ! -f /usr/local/share/ca-certificates/fakeroot.crt ]; do
    wget --timeout 3 --no-check-certificate https://acme-v02.api.letsencrypt.org:15000/roots/0 -O /usr/local/share/ca-certificates/fakeroot.crt
done
update-ca-certificates

if ! ps ax | grep -q "fake_vagrant_sshd[.]py"; then
    python3 /usr/local/sbin/fake_vagrant_sshd.py &>/var/log/fake_vagrant_sshd.log &
fi

exit 0
