# router/rc.local
#!/bin/bash

ip route add 30.0.0.0/16 via 30.0.0.1

sysctl net.ipv4.ip_forward=1
sysctl net.ipv6.conf.all.forwarding=1
